# supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=true

[program:uvicorn]
command=uvicorn main:app --host 0.0.0.0 --port 8081
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/uvicorn_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/uvicorn_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
stdout_events_enabled=true
stderr_events_enabled=true

[program:celery]
command=celery --app=inferflow worker -l INFO --pool solo -Q ${CELERY_QUEUE_NAME}
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/celery_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/celery_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
stdout_events_enabled=true
stderr_events_enabled=true
